# jsk_spot_robot/jsk_spot_startup/services/jsk-spot-utils-reboot.timer
[Unit]
Description=Timer for JSK Spot Reboot Service

[Timer]
OnCalendar = *-*-* 04:00:00

[Install]
WantedBy=timers.target
