# jsk_spot_robot/jsk_spot_startup/services/jsk-spot-utils-auto-service-restart.timer
[Unit]
Description=Timer for JSK Spot Auto Update Workspace

[Timer]
OnCalendar = *-*-* 04:00:00

[Install]
WantedBy=timers.target
