# jsk_spot_robot/jsk_spot_startup/services/jsk-spot-utils-reboot.service
[Unit]
Description=JSK Spot Reboot Service
RefuseManualStart = true
RefuseManualStop = true

[Service]
ExecStart=/sbin/reboot
