# jsk_spot_robot/jsk_spot_startup/services/jsk-spot-utils-network.service
[Unit]
Description=JSK Spot Network Service

[Service]
Type=simple
ExecStart=/var/lib/robot/update-network-connection.sh
ExecStop=/bin/kill -WINCH ${MAINPID}
KillMode=control-group
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
