# jsk_spot_robot/jsk_spot_startup/services/jsk-spot-utils-update-workspace.timer
[Unit]
Description=Timer for JSK Spot Update Workspace Service

[Timer]
OnCalendar = *-*-* 06:00:00

[Install]
WantedBy=timers.target
