# jsk_spot_robot/jsk_spot_startup/services/jsk_spot_networkd.service
[Unit]
Description=JSK Spot Network Service

[Service]
Type=simple
ExecStart=/var/lib/robot/update-network-connection.sh
ExecStop=/bin/kill -WINCH ${MAINPID}
KillSignal=SIGCONT
