# jsk_spot_robot/jsk_spot_startup/config/cockpit.socket.d/listen.conf
[Socket]
ListenStream=
ListenStream=21443
