# jsk_spotkinova_robot/jsk_spotkinova_startup/services/jsk-spotkinova-utils-kinova-rest.timer
[Unit]
Description=Timer for JSK SpotKinova Utils Kinova Rest Service

[Timer]
OnCalendar = *-*-* 07:10:00

[Install]
WantedBy=timers.target
